# scripts/com.hanabi-board.monthly-jouhou.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hanabi-board.monthly-jouhou</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>/Users/yoheimizuno/hanabi-dashboard/scripts/monthly_jouhou.sh</string>
    </array>

    <!-- 毎月1日 9:30 JST 実行 (auto-deployは8:15なので、 後に走らせる) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Day</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/yoheimizuno/hanabi-dashboard/logs/monthly_jouhou.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/yoheimizuno/hanabi-dashboard/logs/monthly_jouhou.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/yoheimizuno</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
